# supervisord.conf
[inet_http_server]
port = 127.0.0.1:9001
username = nanfang
password = nanfang

[supervisord]
logfile = /opt/elogger/logs/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
childlogdir = /opt/elogger/logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://127.0.0.1:9001
username = nanfang
password = nanfang

[program:elogger]
command = python application.py
directory = /opt/elogger/deployed/eventlogger
process_name = %(program_name)s
redirect_stderr = true
stdout_logfile = /opt/elogger/logs/elogger.log
autorestart = true
